# deploy/wm26-backup.timer
[Unit]
Description=WM26 Datenbank-Backup, nächtlich

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true

[Install]
WantedBy=timers.target
